# config/buss2.socket
[Unit]
Description=buss2 socket

[Socket]
ListenStream=/run/buss2.sock
User=caddy

[Install]
WantedBy=sockets.target